# deploy/git-push/crontab
0 11,23 * * * /usr/local/bin/git-push-if-needed.sh >> /proc/1/fd/1 2>&1
